# supervisord.conf
[supervisord]
nodaemon=true

[program:prometheus]
command=/opt/prometheus/prometheus --config.file=/opt/prometheus/prometheus.yml
autostart=true
autorestart=true

[program:grafana]
command=/opt/grafana/bin/grafana-server
autostart=true
autorestart=true

[program:minio]
command=minio server /data --console-address ":9001"
autostart=true
autorestart=true

[program:openmetadata]
command=java -jar /opt/openmetadata/bin/openmetadata.jar server
autostart=true
autorestart=true

[program:airflow]
command=airflow webserver -p 8080
autostart=true
autorestart=true